# pluggramd/pluggramd.service
[Unit]
Description=Built-in program daemon for LED screen web service
Requires=screend.service

[Service]
Type=simple
ExecStart=/home/pi/ledscreen/pluggramd/run.sh
WorkingDirectory=/home/pi/ledscreen/pluggramd/
Restart=always

[Install]
WantedBy=networking.target
